# deploy/launchd/com.mpm.report.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mpm.report.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/phucnt/workspace/my-project-manager/deploy/launchd/run-report.sh</string>
        <string>--daily</string>
    </array>
    <!-- Daily standup digest at 09:00 every day. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/phucnt/workspace/my-project-manager/.data/cron-daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/phucnt/workspace/my-project-manager/.data/cron-daily.err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
